# erjigit17-3420v3/files/etc/init.d/telnetd
#!/bin/sh /etc/rc.common
START=50

start() {
	telnetd -l /bin/login <> /dev/null 2>&1        
	#telnetd -l /usr/libexec/login.sh
}

stop() {
        killall telnetd
}