# erjigit17-3420v3/files/etc/crontabs/root
17 */1 * * * /etc/update_init.sh
*/1 * * * * /etc/check_mccmnc.sh

