# erjigit17/files/etc/uci-defaults/99-4gmodem
#!/bin/sh

ln -s /www/cgi-bin/4gmodem.sh /usr/bin/4gmodem
exit 0
